# scripts/cowork-vm-service.service
[Unit]
Description=Claude Desktop Cowork VM Service Daemon
Documentation=https://github.com/aaddrick/claude-desktop-debian

[Service]
Type=simple
ExecStart=/usr/lib/claude-desktop/node_modules/electron/dist/electron /usr/lib/claude-desktop/node_modules/electron/dist/resources/app.asar.unpacked/cowork-vm-service.js
Environment=ELECTRON_RUN_AS_NODE=1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
